# ubuntu/crontab
# m h  dom mon dow   command
00 * * * * /home/potateros/bin/wallchanger
